# root/etc/s6-overlay/s6-rc.d/init-mod-envarr-postinstall/run
#!/bin/bash

python3 /etc/s6-overlay/s6-rc.d/init-mod-envarr-postinstall/run.py